# logging/crontab
# restart goaccess after a new log file has been started
49,51,53,55,57,59,01 * * * * root killall -TERM goaccess
